# matterbridge/rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Community Hass.io Add-ons: MatterBridge
# Runs the NGINX proxy
# ==============================================================================

bashio::log.info "Starting NGINX..."
exec nginx
