# matterbridge-zigbee2mqtt/rootfs/etc/s6-overlay/s6-rc.d/matterbridge-zigbee2mqtt-installer/up
/etc/s6-overlay/s6-rc.d/matterbridge-zigbee2mqtt-installer/run
