# matterbridge/rootfs/etc/s6-overlay/s6-rc.d/matterbridge-zigbee2mqtt-installer/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Example
# Example init script, runs before any other service
# ==============================================================================

bashio::log.info "Prepare MatterBridge for start"

bashio::log.info "Check MQTT settings"
if ! bashio::config.has_value "host" || ! bashio::config.has_value "port"; then
    bashio::exit.nok "Setting a MQTT host and port is required!"
fi

if ! bashio::config.has_value "topic"; then
    bashio::exit.nok "Setting a Zigbee2MQTT topic is required!"
fi

bashio::log.info "Check Login data"
if ! bashio::config.has_value 'username' || ! bashio::config.has_value 'password'; then
    bashio::exit.nok "Setting a username and password is required!"
fi

bashio::log.info "Linking addon config"
mkdir -p /config/.matterbridge
ln -s /config/.matterbridge /root/.matterbridge

mkdir -p /config/Matterbridge
ln -s /config/Matterbridge /root/Matterbridge

bashio::log.info "Dumping /data/options.json"
echo "$CONFIG_PATH"
cat "$CONFIG_PATH" || true
ls -l /

bashio::log.info "Configuring Zigbee2MQTT plugin"
cat <<EOF | tee /root/.matterbridge/matterbridge-zigbee2mqtt.config.json
{
  "name": "matterbridge-zigbee2mqtt",
  "type": "DynamicPlatform",
  "unregisterOnShutdown": false,
  "host": "$(bashio::config "host")",
  "port": $(bashio::config "port"),
  "topic": "$(bashio::config "topic")",
  "username": "$(bashio::config "username")",
  "password": "$(bashio::config "password")",
  "whiteList": $(bashio::config "whiteList"),
  "blackList": $(bashio::config "blackList")
}
EOF

bashio::log.info "Adding Zigbee2MQTT plugin to MatterBridge"
cd /usr/local/lib/node_modules || exit
npx matterbridge -add matterbridge-zigbee2mqtt
