# matterbridge/rootfs/etc/s6-overlay/s6-rc.d/matterbridge/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: MatterBridge
# Runs matterbridge script
# ==============================================================================

bashio::log.info "Starting Example1..."

exec /usr/bin/matterbridge.sh
